# ansible/roles/system/templates/backup-docker-data.timer.j2
{{ ansible_managed | comment }}
[Unit]
Description=Daily backup of Docker service data

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
